# crontab
0,15,30,45 * * * * python /discursive/index_twitter_stream.py /discursive/topics.txt
